# .config/systemd/user/battery.service
[Unit]
Description=Check Battery Level

[Service]
Type=oneshot
ExecStart=/home/jfa/projects/system-monitor/battery.sh
