# .config/stumpwm.d/stumpwm-session.target
[Unit]
Description=Xorg Session Running StumpWM

[Service]
BindsTo=graphical-session.target
